# etc/systemd/system/bluebinder.service.d/10-wait.conf
[Service]
ExecStartPre=
ExecStartPre=/usr/bin/waitforservice init.svc.*synabtlinux*
