# etc/systemd/system/android-service@hwcomposer.service.d/30-no-restart.conf
[Service]
ExecStart=
ExecStop=
ExecStart=/bin/true
ExecStop=/bin/true
