# etc/systemd/system/droidian-fpd-client.service.d/10-wait.conf
[Service]
ExecStartPre=
ExecStartPre=/usr/bin/binder-wait -d /dev/binder android.hardware.biometrics.fingerprint.IFingerprint/default
