# etc/systemd/system/android-mount.service.d/50-no-phosh.conf
[Service]
ExecStartPre=-/usr/bin/rm -f /etc/systemd/system/phosh.service.d/40-droid-wait.conf
ExecStartPre=-/usr/bin/systemctl daemon-reload
ExecStartPre=-/usr/bin/systemctl enable phosh
